# crontab
*/15 * * * * $BACKUP_WORK_DIR/backup.sh 2>&1
